# CyberStakes_2020/turtles_all_the_way_down/controlled_access/_firmware.bin.extracted/squashfs-root/usr/lib/opkg/info/at.list
/usr/bin/atrm
/etc/init.d/atd
/usr/bin/atq
/usr/bin/at
/usr/sbin/atd
