# cyberstakes/2020/controlled_access/_firmware.bin.extracted/squashfs-root/usr/lib/opkg/info/autossh.list
/etc/hotplug.d/iface/20-autossh
/etc/config/autossh
/etc/init.d/autossh
/usr/sbin/autossh
