# cyberstakes/2020/controlled_access/_firmware.bin.extracted/squashfs-root/usr/lib/opkg/info/openssh-server.list
/usr/sbin/sshd
/etc/init.d/sshd
/etc/ssh/sshd_config
/lib/upgrade/keep.d/openssh-server
